# board/demo/rootfs_overlay/etc/profile.d/50-bconenv.sh
export BCON_ADAPTER_I2C_DEVICES="/dev/i2c-0:60"