# package/libbconctl/50-gpioled.rules
DEVPATH=="/devices/soc0/gpio-leds/leds/led-user[012]", RUN+="/usr/lib/bconctl/setup_leds"
